# services/jukebox-frontend.service
[Unit]
Description=Jukebox Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=jukebox
Group=jukebox
WorkingDirectory=/Users/sloor/Developer/jukebox_frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PATH=/usr/bin:/usr/local/bin:/bin:/usr/sbin:/usr/local/sbin

ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target